# cookbooks/jdabbs/templates/default/vlc.upstart.conf.erb
description "vlc http daemon"

start on runlevel [2345]
stop on runlevel [016]

expect daemon

setuid james

exec /usr/bin/vlc -I http --daemon
